# deploy/ansible/roles/forge-agent/templates/forge-agent.service.j2
# Managed by Ansible -- do not edit manually
[Unit]
Description=Forge Agent - Docker sandbox orchestration agent
Documentation=https://github.com/appx/forge
After=network-online.target docker.service
Wants=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/forge-agent
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
EnvironmentFile=/etc/forge/agent.env

# Security hardening (T-02-18 mitigation)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/forge
PrivateTmp=true

# File descriptor limit for many concurrent containers
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forge-agent

[Install]
WantedBy=multi-user.target
